# config/includes.chroot/etc/skel/.config/autostart/neuros-tray.desktop
[Desktop Entry]
Type=Application
Name=NeurOS Tray
Comment=NeurOS System Tray Applet
Exec=/usr/local/bin/neuros-tray
Icon=neuros
Terminal=false
Categories=Utility;
X-GNOME-Autostart-enabled=true
